# client/usb-monitor.service
[Unit]
Description=USB Monitor Security Daemon
After=network.target systemd-logind.service
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/usb-monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usb-monitor

# Безопасность (ослаблены ограничения для монтирования и уведомлений)
NoNewPrivileges=false
ProtectHome=false
ReadWritePaths=/media /var/log /tmp /proc /sys /dev /run /etc
PrivateTmp=false
PrivateMounts=false
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=false
RestrictRealtime=true
RestrictNamespaces=false
ProtectSystem=false
PrivateDevices=false

# Переменные окружения
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
